# build/etc/profile.d/docker-env.sh
export DOCKER_HOST=$(ip route ls default | gawk '{print $3}'):23900