# build/etc/profile.d/python-env.sh
export HTTPLIB2_CA_CERTS=/etc/ssl/cert.pem
export REQUESTS_CA_BUNDLE=/etc/ssl/cert.pem